# dist/zero_os_bundle_20260308T210753Z/build/native_store_prod/backend_deploy/zero-store-backend.service
[Unit]
Description=Zero Store Backend
After=network.target

[Service]
WorkingDirectory=/opt/zero-store
ExecStart=/usr/bin/python3 -c 'from zero_os.native_store_backend import run_server; run_server("/opt/zero-store")'
Restart=always

[Install]
WantedBy=multi-user.target
